# deploy/uberspace/nuola.ini
; supervisord-Dienste fuer das Nuola Energy Dashboard auf Uberspace 8.
;
; Ablage auf dem Host:  ~/etc/services.d/nuola.ini
; Aktivieren:           supervisorctl reread && supervisorctl update
; Status:               supervisorctl status nuola-web nuola-worker
;
; Vor dem Kopieren die drei Platzhalter ersetzen:
;   USERNAME  -> der Uberspace-Benutzername (z.B. isabell)
;   PORT      -> ein freier Port 1024-65535, identisch zum Web-Backend
;                (uberspace web backend set / --http --port PORT)
;
; Beide Dienste laden die Umgebungsvariablen aus der .env im Projektordner
; (per `source`), damit der Worker (reiner tsx-Prozess) dieselben Variablen wie
; der Next.js-Server bekommt. Details: deploy/uberspace/README.md

[program:nuola-web]
directory=/home/USERNAME/nuola
command=bash -c 'set -a; source /home/USERNAME/nuola/.env; set +a; exec /home/USERNAME/nuola/node_modules/.bin/next start -H 0.0.0.0 -p PORT'
autostart=true
autorestart=true
startsecs=30
stopwaitsecs=20
stdout_logfile=/home/USERNAME/logs/nuola-web.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
redirect_stderr=true
environment=NODE_ENV="production",HOME="/home/USERNAME"

[program:nuola-worker]
directory=/home/USERNAME/nuola
command=bash -c 'set -a; source /home/USERNAME/nuola/.env; set +a; exec /home/USERNAME/nuola/node_modules/.bin/tsx src/worker/index.ts'
autostart=true
autorestart=true
startsecs=30
stopwaitsecs=20
stdout_logfile=/home/USERNAME/logs/nuola-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
redirect_stderr=true
environment=NODE_ENV="production",HOME="/home/USERNAME"
